# localai-config/com.localai.service.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.localai.service</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Users/pernelltoney/LocalAI/local-ai</string>
        <string>--models-path=/Users/pernelltoney/LocalAI/models</string>
        <string>--address=:8080</string>
        <string>--debug=false</string>
    </array>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <true/>
    
    <key>StandardOutPath</key>
    <string>/Users/pernelltoney/LocalAI/logs/localai.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/pernelltoney/LocalAI/logs/localai-error.log</string>
    
    <key>WorkingDirectory</key>
    <string>/Users/pernelltoney/LocalAI</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
